# map_directory/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "/root/timerprint.sh"
